# speech2text.service
[Unit]
Description=Speech-to-Text Global Hotkey Daemon
After=network.target

[Service]
Type=simple
User=%u
ExecStart=/usr/bin/python3 %h/.local/share/speech-to-text-hotkey/speech2text.py
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
Environment="STT_HOTKEY=z"

[Install]
WantedBy=default.target
