# typefree.service
[Unit]
Description=Typefree - Speech-to-Text global hotkey daemon
# Bind to the graphical session: PartOf + WantedBy=graphical-session.target make
# the After= ordering effective, so the daemon only starts once the compositor
# has exported WAYLAND_DISPLAY/DISPLAY. Without this the tray's Qt init races the
# session at boot and the process core-dumps (qFatal: no display).
After=graphical-session.target ydotoold.service
PartOf=graphical-session.target
Wants=ydotoold.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.local/share/typefree/typefree.py
Restart=on-failure
RestartSec=3
Environment=TYPEFREE_HOTKEY=z
Environment=TYPEFREE_MODIFIER=alt
Environment=TYPEFREE_MODE=hold
Environment=YDOTOOL_SOCKET=%t/.ydotool_socket
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical-session.target
